# koha/rootfs/etc/s6-overlay/s6-rc.d/zebra-server/run
#!/command/execlineb -P
with-contenv

# Redirect stderr to stdout for logging
fdmove -c 2 1

importas KOHA_INSTANCE KOHA_INSTANCE

# Define and export essential Koha environment variables
export PERL5LIB /usr/share/koha/lib
export KOHA_CONF /etc/koha/sites/${KOHA_INSTANCE}/koha-conf.xml
export KOHA_HOME /usr/share/koha

koha-shell ${KOHA_INSTANCE} -c \
  "/usr/bin/zebrasrv -v 1 -f /etc/koha/sites/${KOHA_INSTANCE}/koha-conf.xml"
